# encfs-db-corruption-monitor.service
[Unit]
Description=Encfs/Dropbox corruption monitor

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=COMMAND ARGS
Type=simple
